# services/typotuner.service
[Unit]
Description=TypoTuner Typing Analysis Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/typotuner start -f
Restart=on-failure
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
